# docker/tornado.ini
[program:tornado]
command=python server.py --port=$PORT
directory=/app
#user=root
autorestart=true
redirect_stderr=true
stdout_logfile=logs/supervisor.log
loglevel=info
